# chimetime.service
[Unit]
Description=Chime Time
After=multi-user.target

[Service]
User=root
WorkingDirectory=/root/chimetime
ExecStart=/root/chimetime/venv/bin/python ct_webserver.py

Restart=always

[Install]
WantedBy=multi-user.target
