# chimetime.service
[Unit]
Description=Chime Time
After=multi-user.target

[Service]
User=pi
WorkingDirectory=/home/pi/chimetime
ExecStart=/usr/bin/python3 -m gunicorn --bind=0.0.0.0:5000 "ct_webserver:app_factory('config.json', fake=False)"

Restart=always

[Install]
WantedBy=multi-user.target
